# src/initramfs/inittab
::sysinit:/etc/init.d/rcS
tty1::respawn:/sbin/getty 38400 tty1
